# bin/ord_bot.service
[Unit]
After=network.target
Description=ord_bot
StartLimitBurst=120
StartLimitIntervalSec=10m

[Service]
User=root
WorkingDirectory=/var/lib/ord_bot
ExecStart=/usr/local/bin/pipenv run python3.10 /var/lib/ord_bot/ord_bot.py
Restart=on-failure
RestartSec=5s
StandardOutput=syslog+console
StandardError=syslog+console
TimeoutStopSec=10m
Type=simple

[Install]
WantedBy=multi-user.target
